# deploy/midnight.service
# Ubuntu 24.04 systemd: /etc/systemd/system/midnight.service
# sudo systemctl daemon-reload && sudo systemctl enable --now midnight
#
# Adjust User/Group and paths if your app lives somewhere other than /srv/midnight

[Unit]
Description=Midnight (Gunicorn)
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/midnight
EnvironmentFile=/srv/midnight/.env
ExecStart=/srv/midnight/.venv/bin/gunicorn \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --timeout 60 \
    midnight.wsgi:application
Restart=on-failure

[Install]
WantedBy=multi-user.target
